# init.d/cron
#!/sbin/runscript

command=cron
command_args=-d
name="cron"
