# init.d/crond
#!/sbin/runscript

command=crond
command_args=-d
name="cron daemon"
