# init.d/dcron
#!/sbin/runscript

command=/usr/bin/dcron
command_args=-d
name="dcron cron daemon"
